# deploy/nexus.service
[Unit]
Description=Agent Nexus Dashboard
After=network.target

[Service]
Type=simple
User=pets
WorkingDirectory=/home/pets/zoo/agents_sessions_dashboard
Environment="PYTHONPATH=/home/pets/zoo/agents_sessions_dashboard"
EnvironmentFile=/home/pets/zoo/agents_sessions_dashboard/.env
ExecStart=/bin/sh -lc '/usr/bin/python3 -m uvicorn backend.api.main:app --host "${NEXUS_HOST:-0.0.0.0}" --port "${NEXUS_BACKEND_PORT:-18890}"'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
